# build/config/includes.chroot/etc/skel/.config/autostart/tahoe-session.desktop
[Desktop Entry]
Type=Application
Name=Tahoe Session Bootstrap
Exec=/usr/local/bin/tahoe-session-start
X-GNOME-Autostart-enabled=true
NoDisplay=false
Terminal=false
