# build/config/includes.chroot/etc/skel/.config/autostart/tahoe-startup-animation.desktop
[Desktop Entry]
Type=Application
Name=Tahoe Startup Animation
Exec=/usr/local/bin/tahoe-startup-animation
X-GNOME-Autostart-enabled=true
Terminal=false
