# build/config/includes.chroot/etc/skel/.config/autostart/tahoe-dynamic-wallpaper.desktop
[Desktop Entry]
Type=Application
Name=Tahoe Dynamic Wallpaper
Exec=/usr/local/bin/tahoe-dynamic-wallpaper-daemon
X-GNOME-Autostart-enabled=true
Terminal=false
