# build/config/includes.chroot/etc/skel/.config/autostart/tahoe-setup-assistant.desktop
[Desktop Entry]
Type=Application
Name=Tahoe Setup Assistant
Exec=/usr/local/bin/tahoe-setup-assistant
X-GNOME-Autostart-enabled=true
Terminal=false
